# deploy/nowenos.service
[Unit]
Description=NoWenOS NAS Management Panel
After=network.target docker.service
Wants=docker.service

[Service]
Type=simple
User=nowenos
Group=nowenos
WorkingDirectory=/opt/nowenos
EnvironmentFile=/etc/nowenos/nowenos.env
ExecStart=/opt/nowenos/nowenos-api
Restart=always
RestartSec=5
StandardOutput=append:/var/log/nowenos/nowenos.log
StandardError=append:/var/log/nowenos/nowenos.log

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/var/lib/nowenos /var/log/nowenos
ProtectHome=true

[Install]
WantedBy=multi-user.target
